# root/etc/services.d/web/run
#!/usr/bin/with-contenv bash

# Run App in development mode
cd /app/ffmpeg-web
exec \
	s6-setuidgid abc nodemon \
	--exec python3 server.py > /applogs/app.log 2>&1
